# dot_config/systemd/user/cc-inflight-reminder.service
[Unit]
Description=Scan for in-flight Claude Code work and notify

[Service]
Type=oneshot
ExecStart=%h/.claude/hooks/daily-inflight-reminder.sh
Environment=DISPLAY=:0
