# dot_config/systemd/user/cc-inflight-reminder.timer
[Unit]
Description=Daily reminder for in-flight Claude Code work

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
